# systemd/loopback-port-update.service
[Unit]
Description=Update loopback target port from VPN port forwarding

[Service]
Type=oneshot
ExecStart=%h/loopback/update-loopback-port.sh
